# usbmount.rules
ACTION=="add", KERNEL=="sd*[!0-9]", RUN+="/usr/local/bin/usbmount.sh %N"
