# usbmount.rules
ACTION=="add", KERNEL=="sd*[!0-9]", RUN+="/bin/systemctl --no-block start piban@%k.service"
